# steps/radiodan-updater/files/install-updates
#!/bin/sh
### BEGIN INIT INFO
# Provides:          install-updates
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:
# Short-Description: Installs pending updates
# Description:       Installs any pending radiodan app updates
### END INIT INFO

set -e

export PATH=$PATH:/usr/local/bin

# Install any pending services (if any)
/usr/local/bin/updater install --workspace=/opt/radiodan/updates --target=/opt/radiodan/apps &