# config/radiodan-physical-ui.conf
start on runlevel [2345]
stop on runlevel [016]

respawn

setuid root
setgid root

umask 0002

env HOME=/root

chdir /opt/radiodan/physical-ui

exec bash -lc './bin/ui'
