# config/wpa-conf-copier.conf
author "Andrew Nicolaou <andrew.nicolaou@bbc.co.uk>"
description "Copy wpa_supplicant config from boot partition on startup"

task

start on starting networking

env SOURCE_CONF="/boot/wpa_supplicant.txt"
env TARGET_CONF="/etc/wpa_supplicant/wpa_supplicant.conf"

script
  if [ "$SOURCE_CONF" -nt "$TARGET_CONF" ]; then
    echo "$SOURCE_CONF newer than $TARGET_CONF - copied"
    /bin/cp -u $SOURCE_CONF $TARGET_CONF
  else
    echo "$SOURCE_CONF older than $TARGET_CONF - not copied"
  fi
  exit 0
end script
